# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:msdApp]
user=root
directory=/opt/msd
stdout_logfile=/var/log/msd.log
command=/usr/bin/gulp run