# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:msdApp]
command=npm start /src > /logs/msd_log.txt -D